# private_geth.service
[Service]
ExecStart=/home/ubuntu/ethereum/geth_scripts/run.sh

[Install]
WantedBy=default.target
